# systemd/hermes-xte-kitchen-server.service
[Unit]
Description=Hermes: xte-kitchen-server (Xteink X4 kitchen-display host)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=bot
Group=bot
WorkingDirectory=/home/bot/apps/xte-kitchen-server
Environment=XTE_KITCHEN_STATE_DIR=/home/bot/apps/xte-kitchen-server/state
Environment=XTE_KITCHEN_SECRETS_DIR=/home/bot/apps/xte-kitchen-server/secrets
Environment=XTE_KITCHEN_BIND_HOST=0.0.0.0
Environment=XTE_KITCHEN_BIND_PORT=8080
ExecStart=/home/bot/.local/bin/xte-kitchen-server
Restart=on-failure
RestartSec=5

# Best-effort hardening (LAN-only service; not security-critical):
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=read-only
PrivateTmp=true

[Install]
WantedBy=multi-user.target
